# applications/globaleaks/cron/crontab
#crontab
*      *  *  *  *  root *applications/globaleaks/cron/scheduler.py
0      */2  *  *  *  root *applications/globaleaks/cron/target_list_clean.py
