# applications/globaleaks/cron/crontab
#crontab
*      *  *  *  *  root *applications/globaleaks/cron/scheduler.py
#XXX To be removed in non demo usage
0      */2  *  *  *  root *applications/globaleaks/cron/target_list_clean.py
