# globaleaks/applications/globaleaks/cron/crontab
# crontab file format:
#
# remind - need to be checked - we are not really run as root ?
*      *  *  *  *  root *applications/globaleaks/cron/scheduler.py
